# boulderryn.service
[Unit]
Description=BoulderRyn Climbing Gym Management Server
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/.openclaw/workspace/boulderryn-project
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=5
Environment=PORT=8080
Environment=NODE_ENV=production
StandardOutput=journal
StandardError=journal
SyslogIdentifier=boulderryn

[Install]
WantedBy=multi-user.target
